# server/etc/systemd/system/bnl-discord-bot.service
[Unit]
Description=Start the app

[Service]
EnvironmentFile=/var/live/bnls/.env
ExecStart=/usr/bin/node /var/live/bnls/bot.bin.js
Restart=always
