# server/etc/systemd/system/dnl-discord-bot.service
[Unit]
Description=Start the app

[Service]
ExecStart=/usr/bin/node /var/live/bnl-discord-bot/bot.bin.js
Restart=always
